# execution/mesh/proxy/systemd/mesh-signal-worker.service
[Unit]
Description=Mesh Signal Worker (stdio jsonRpc)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=signal
Group=signal
WorkingDirectory=/opt/Joi/execution/mesh/proxy
EnvironmentFile=-/etc/default/mesh-signal-worker
ExecStart=/opt/Joi/execution/mesh/proxy/run-worker.sh
Restart=always
RestartSec=3

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/signal-cli
RuntimeDirectory=signal-cli
RuntimeDirectoryMode=0755
UMask=0077

[Install]
WantedBy=multi-user.target
